# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/eb_docker_study/bin/uwsgi --ini /srv/eb_docker_study/.config/uwsgi/uwsgi-app.ini